# files/2_monitor/supervisor/DSM101.conf
[program:run_DSM101]
command = /opt/monitor/script/DSM101.py
process_name = %(program_name)s
autostart = false
autorestart = true
user = pi
redirect_stderr = true
stdout_logfile = /opt/monitor/log/DSM101.out
stderr_logfile = /opt/monitor/log/DSM101.err
environment = DEVICE_MD = "DSM101", DEVICE_USB = "/dev/serial/by-id/usb-Prolific_Technology_Inc._USB-Serial_Controller_D-if00-port0", DEVICE_POS = "AMoRE_EXP_HALL"

